# onetouch-autostart-daemon.rules
ATTRS{idVendor}=="0d49", ATTRS{idProduct}=="7201", ENV{DEVTYPE}=="disk", RUN+="onetouch-iii-button '%r/$name'"
